# ssh/supervisord-cron.conf
[program:crond]
command=/usr/sbin/crond -f
autorestart=true
startretries=3
stdout_logfile=/var/log/crond.log
