# base/supervisor.d/crond.ini
[program:crond]
command=crond -f
priority=999
autostart=true
startsecs=1
autorestart=true
user=root

