# base/supervisor.d/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D -e
priority=999
autostart=true
startsecs=1
autorestart=true
user=root

