# base/supervisor.d/syslog.ini
[program:syslog]
command=syslog-ng -F
priority=999
autostart=true
startsecs=1
autorestart=true
user=root

