# php5/supervisord-php5.conf
[program:php-fpm]
command=php-fpm5 -F
stdout_logfile=/var/log/php.log
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
startretries=50

