# java/supervisord-java.conf
[program:jar]
command=java -jar /data/web/web.jar
stdout_logfile=/var/log/web.log
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
startretries=50

